# crontab
* * * * * cd /var/www/html && php bin/console app:check-pomodoro-sessions >> /var/log/cron.log 2>&2